# deploy/systemd/xin-channel-gateway.service
[Unit]
Description=Xin Channel Gateway
After=xin-orchestrator.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=xin
Group=xin
WorkingDirectory=/opt/xin-chatbot/src
EnvironmentFile=/opt/xin-chatbot/config/gateway.env
Environment="PATH=/opt/xin-chatbot/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin"
Environment="PYTHONPATH=/opt/xin-chatbot/src"
ExecStart=/usr/bin/env bash -c 'source /opt/xin-chatbot/.venv/bin/activate && exec uvicorn chatbot.apps.gateway.app:create_app --factory --host "${GATEWAY_BIND:-127.0.0.1}" --port "${GATEWAY_PORT:-8080}"'
ExecStartPost=/usr/bin/env bash -c '/opt/xin-chatbot/src/scripts/systemd/healthcheck.sh --url "http://127.0.0.1:${GATEWAY_PORT:-8080}/health" --name xin-channel-gateway'
Restart=on-failure
RestartSec=5
TimeoutStartSec=60
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
